# systemd/weatherhat.service
[Unit]
Description=WeatherHAT Monitoring Service
Documentation=https://github.com/faintdeception/weather-station
After=network.target mongod.service
Wants=mongod.service

[Service]
Type=simple
User=REPLACE_USER
Group=REPLACE_USER
WorkingDirectory=REPLACE_REPO_DIR
Environment="PATH=REPLACE_HOME/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=REPLACE_REPO_DIR"
Environment="WEATHER_LOCATION=backyard"
Environment="WEATHER_INTERVAL=60"
EnvironmentFile=-REPLACE_REPO_DIR/.env
ExecStart=/usr/bin/python3 REPLACE_REPO_DIR/weatherhat_service.py
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=300

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=weatherhat

# Security
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=REPLACE_REPO_DIR

[Install]
WantedBy=multi-user.target
